# src/deb/nginx/postrm
#!/bin/sh

update-rc.d hestia remove >/dev/null 2>&1


exit 0
