# src/deb/nginx/postrm
#!/bin/sh

update-rc.d devcp remove > /dev/null 2>&1

exit 0
