# src/deb/web-terminal/hestia-web-terminal.service
[Unit]
Description=HestiaCP Web Terminal
Documentation=https://hestiacp.com/docs/
After=network.target

[Service]
User=root
Environment=NODE_ENV=production
Environment=HESTIA=/usr/local/hestia
ExecStart=$HESTIA/web-terminal/server.js
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
